# projects/alpaca-trading/config/daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.alpaca-trader</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/raythomas/.openclaw/workspace/projects/alpaca-trading/.venv/bin/python</string>
        <string>/Users/raythomas/.openclaw/workspace/projects/alpaca-trading/src/main.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/raythomas/.openclaw/workspace/projects/alpaca-trading</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/raythomas/.openclaw/workspace/projects/alpaca-trading</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/raythomas/.openclaw/workspace/projects/alpaca-trading/logs/stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/raythomas/.openclaw/workspace/projects/alpaca-trading/logs/stderr.log</string>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
